# crx_download.crontab
SHELL=/bin/bash
PATH=/home/ubuntu/dbling/src:/usr/bin:/bin
MAILTO=nobody@example.com
HOME=/home/ubuntu/dbling/src
XDG_RUNTIME_DIR=/run/user/1000

42 9,21 * * * { time /home/ubuntu/dbling/src/crx_crawler.py -up; } 2>&1 | mail -s "dbling: Result of list update" $MAILTO
47 9,21 * * * { time /home/ubuntu/dbling/src/crx_crawler.py -sp; echo -e "\n"; cat /home/ubuntu/dbling/src/stats.json; } 2>&1 | mail -s "dbling: Result of centroid calculations" $MAILTO
#32 9 * * 0 /usr/bin/sudo /sbin/reboot

